# systemd/dashboard.service
# User service: ~/.config/systemd/user/dashboard.service
# Enable with:  systemctl --user enable --now dashboard
# Serves on 0.0.0.0:8090; expose to your phone with: tailscale serve --bg 8090
[Unit]
Description=Logos Node Dashboard
After=network-online.target logos-node.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/logos-blockchain-runbook
Environment=HOST=0.0.0.0
Environment=PORT=8090
# ZONE_CHANNEL: the operator's own zone channel name (his handle on the network)
Environment=ZONE_CHANNEL=optiplex
# Wallet key is auto-derived from configs/live/node.yaml (funding_pk) by dashboard/run.sh
ExecStart=%h/logos-blockchain-runbook/dashboard/run.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
